# kb-extra-phone/files/etc/uci-defaults/71_kb_phone
#!/bin/sh

. /lib/functions.sh

### Enable asterisk
uci -q batch <<-EOF
	set asterisk.general.enabled='1'
EOF
uci commit
exit 0
