# kb-extra-app/files/etc/crontabs/root
*/1 * * * * /opt/kb-script/update_redir.sh 
